# container/supervisord.conf
[supervisord]
nodaemon=true
user=hermes
logfile=/opt/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:desktop]
command=/opt/hermes-box/start-desktop.sh
priority=10
autorestart=true
startsecs=5
stdout_logfile=/opt/data/logs/desktop.log
stderr_logfile=/opt/data/logs/desktop.err.log

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5901
priority=20
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/opt/data/logs/novnc.log
stderr_logfile=/opt/data/logs/novnc.err.log

[program:keeper]
command=/bin/bash -lc 'while sleep 60; do used=$(du -sk /opt/data | cut -f1); budget=$(( ${DISK_BUDGET_GB:-10} * 1024 * 1024 )); if (( used > budget )); then echo "WARNING data exceeds soft disk budget: ${used}KiB > ${budget}KiB"; fi; done'
priority=30
autorestart=true
stdout_logfile=/opt/data/logs/disk-budget.log
stderr_logfile=/opt/data/logs/disk-budget.err.log
